# gnome15/src/gnome/autostart/gnome15.desktop.in
[Desktop Entry]
Version=1.0
Encoding=UTF-8
Name=Logitech G Keyboard Desktop Service
Icon=gnome15
Comment=Logitech G series keyboard desktop service.
Exec=g15-desktop-service
Terminal=false
Type=Application
Categories=
GenericName=
#X-GNOME-Autostart-Delay=8
#AutostartCondition=GNOME /desktop/gnome/gnome15/enabled
#X-GNOME-Autostart-Phase=Applications
#X-GNOME-AutoRestart=true
#X-Ubuntu-Gettext-Domain=gnome15